# TP3/youtube.service
[Unit]
Description="Youtube Download Service"

[Service]
ExecStart=sudo bash /srv/youtube/youtubev2.sh

[Install]
WantedBy=multi-user.target

