# TP3/youtube.service
[Unit]
Description="Youtube DL Service"

[Service]
ExecStart=sudo bash /srv/yt/yt-v2.sh

[Install]
WantedBy=multi-user.target

